# squashfs-root/usr/etc/inittab
::sysinit:/usr/etc/rcS
ttyS1::respawn:/bin/sh
